# add/upstart_comcon-i/nia.conf
# upstart service file at /etc/init/nia.conf
description "NIA Service"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [!2345]


# Automatically restart process if crashed
respawn
respawn limit 10 5
#exec fork
# we don't use buil-in log because we use a script below
# console log

# drop root proviliges and switch to comses user
setuid diver
setgid diver




env HOME="/home/diver"
env USER="diver"


#chdir /home/diver/Apps/gsongs/gsongs

#script
    #env HOME=/home/comcon
exec /home/diver/bin/python3 /home/diver/Apps/gsongs/gsongs/run.py >> /home/diver/ggg.log
#end script
